# T-LAB/roles/lxd-c4-provision/files/prometheus.service
[Service]
User=prometheus
Group=prometheus
ExecStart=/mnt/vars/prometheus/prometheus \
              --config.file=/mnt/vars/prometheus/prometheus.yaml \
              --storage.tsdb.path=/mnt/vars/prometheus/data
ExecReload=/bin/kill -HUP $MAINPID

[Install]
WantedBy=multi-user.target

